# brightness/ddcci-add-device.service
[Unit]
Description=DDCCI Shell Script

[Service]
ExecStart=/usr/bin/ddcci_add_devices.sh

[Install]
WantedBy=multi-user.target
